# tcp_server_logger/rootfs/etc/services.d/tcp_server/run
#!/usr/bin/with-contenv sh
exec python3 /usr/src/app/tcp_server.py
